# packages/kademar-skel/kademar-skel/usr/share/kademar/skel/.config/autostart/helioxhelper.desktop
[Desktop Entry]
Categories=QT;KDE;GNOME;Utility;Accessibility;
Exec=HelioxHelper
GenericName=Heliox Helper
Icon=/usr/share/helioxhelper/icon.png
MimeType=
Name=Heliox Helper
StartupNotify=true
Terminal=false
TerminalOptions=
Type=Application
