# packages/kademar-base/kademar-base/usr/lib/systemd/system/kademar-services.service
[unit]
Description="Kademar Services (old rc.local)"
After=network.target

[Service]
Type=oneshot
ExecStart=/etc/rc.d/kademar-services &
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

