# packages/kademar-base/kademar-base-5.7.1/usr/lib/udev/rules.d/99-kademar-cups.rules
ACTION=="add", KERNEL=="printer", RUN+="/bin/sh -c '[ -x /etc/rc.d/cupsd -a ! -e /var/run/daemons/cupsd ] && rc.d start cupsd >/dev/null'"
ACTION=="add", KERNEL=="lp*", RUN+="/bin/sh -c '[ -x /etc/rc.d/cupsd -a ! -e /var/run/daemons/cupsd ] && rc.d start cupsd >/dev/null'"