# Debian-OLD/kademarcenter/cfg/crontab
0 */6     * * *   root  sh /usr/share/kademar/utils/kademarcenter/scripts/warn_upgrade
