# packages/kademar-base/kademar-base/etc/cron.daily/remove-pacman-cache
#!/bin/bash

#Clean pacman downloaded files

find /var/cache/pacman/pkg  -type f   | xargs rm -f
