# packages/heliox-base/heliox-base/usr/lib/systemd/system/kademar-festival.service
[Unit]
Description=Festival Speech Synthesis System
After=sound.target
Before=multi-user.target

[Service]
Type=oneshot
ExecStart=/etc/rc.d/festivald start
ExecStop=/etc/rc.d/festivald stop
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
