# packages/kademar-livecd/kademar-livecd/usr/lib/systemd/system/kademar-livecd.service
[unit]
Description="Kademar Load HDD Swap Partitions and livecd workarround"
After=basic.target
Before=multi-user.target

[Service]
Type=oneshot
ExecStart=/etc/rc.d/kademar-livecd start &
ExecStop=/sbin/swapoff -a
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

