# packages/kademar-livecd/kademar-livecd/usr/lib/systemd/system/kademar-swap.service
[unit]
Description="Kademar Load HDD Swap Partitions"
After=network.target

[Service]
Type=oneshot
ExecStart=/etc/rc.d/kademar-swap start &

[Install]
WantedBy=multi-user.target

