# packages/kademar-skel/kademar-skel-5.1.3/usr/share/kademar/skel/.config/autostart/simond.desktop
[Desktop Entry]
Encoding=UTF-8
Version=1.0
Name=simond
Type=Application
Exec=simond
X-GNOME-Autostart-enabled=true
Icons=/usr/share/icons/hicolor/48x48/apps/sam.png