# packages/kademar-base/kademar-base/usr/lib/udev/rules.d/99-kademar-cups.rules
ACTION=="add", KERNEL=="printer", RUN+="/bin/sh -c '[ -e /usr/lib/systemd/system/cups.service -a -z "$(systemctl status cups | grep -i running)" ] && systemctl start cups >/dev/null'"
ACTION=="add", KERNEL=="lp*", RUN+="/bin/sh -c '[ -e /usr/lib/systemd/system/cups.service -a -z "$(systemctl status cups | grep -i running)" ] && systemctl start cups >/dev/null'"