# packages/kademar-base/kademar-base/usr/lib/udev/rules.d/99-kademar-bluetooth.rules
ACTION=="add", KERNEL=="bluetooth", RUN+="/bin/sh -c '[ -e /usr/lib/systemd/system/bluetooth.service -a -z "$(systemctl status bluetooth | grep -i running)" ] && systemctl start bluetooth >/dev/null'", OPTIONS+="ignore_device"
ACTION=="add", KERNEL=="hci*", SUBSYSTEM=="bluetooth", RUN+="/bin/sh -c '[ -e /usr/lib/systemd/system/bluetooth.service -a -z "$(systemctl status bluetooth | grep -i running)" ] && systemctl start bluetooth >/dev/null'", OPTIONS+="ignore_device"