# packages/kademar-base/kademar-base/usr/lib/udev/rules.d/99-kademar-bluetooth.rules
ACTION=="add", KERNEL=="bluetooth", RUN+="/bin/sh -c '[ -x /etc/rc.d/bluetooth -a ! -e /var/run/daemons/bluetooth ] && rc.d start bluetooth >/dev/null'", OPTIONS+="ignore_device"
ACTION=="add", KERNEL=="hci*", SUBSYSTEM=="bluetooth", RUN+="/bin/sh -c '[ -x /etc/rc.d/bluetooth -a ! -e /var/run/daemons/bluetooth ] && rc.d start bluetooth >/dev/null'", OPTIONS+="ignore_device"